# profile.d/01-whoami.sh
export WHOAMI="`whoami`"
[ -z "$MYLOGNAME" ] && export MYLOGNAME="$WHOAMI"
