# profile.d/01-whoami.sh
### direct2cache head START ###
export WHOAMI="`whoami`"
[ -z "$MYLOGNAME" ] && export MYLOGNAME="$WHOAMI"
#### direct2cache head END ####
