# profile.d/09-uptime.sh
havebin uptime && {
  case $0 in -*)
    printf "uptime:"
    uptime;;
  esac
}
