# profile.d/09-uptime.sh
have uptime && \
  printf 'uptime: ' && \
  uptime
