# profile.d/09-uptime.sh
# Used as Fish scripts too.

have uptime && \
  printf 'uptime: ' && \
  uptime
