# profile.d/10-iscygwin.sh
is_cygwin() {
  startswith "$TERM" cygwin
}
