# profile.d/6-editor.sh
export EDITOR="$editor"
_echo2cache "export EDITOR='$editor'"
