# profile.d/6-lesspipe.sh
have lesspipe.sh && export LESSOPEN="|lesspipe.sh %s"
