# profile.d/6-lesspipe.sh
[ -x /usr/bin/lesspipe ] && _echo2cache "$(SHELL=/bin/sh lesspipe)"
