# profile.d/6-timestyle.sh
_echo2cache "export TIME_STYLE='$time_style'"
