# profile.d/70-unalias-all.sh
### direct2cache head START ###
unalias -a
#### direct2cache head END ####
