# profile.d/75-alias.sh
### direct2cache head START ###
alias  vialiases='vi "$HOME/.profile.local"'
alias applyaliases='. $HOME/.profile.local; true'

alias        cd-='cd -'
alias        cd.='cd .'
alias       cd..='cd ..'
alias      cd...='cd ../..'
alias     cd....='cd ../../..'
alias    cd.....='cd ../../../..'
alias   cd......='cd ../../../../..'
alias  cd.......='cd ../../../../../..'
alias cd........='cd ../../../../../../..'

alias  langposix="langset POSIX      POSIX"
alias     langzh="langset zh_CN      zh_CN:zh"
alias langzhutf8="langset zh_CN.utf8 zh_CN:zh"

alias proxyclear='unset http_proxy https_proxy ftp_proxy'

#### direct2cache head END ####
