# profile.d/75-langset.sh
### direct2cache head START ###
alias  langposix="langset POSIX      POSIX"
alias     langzh="langset zh_CN      zh_CN:zh"
alias langzhutf8="langset zh_CN.utf8 zh_CN:zh"
#### direct2cache head END ####
