# profile.d/75-langset.sh
### direct2cache head START ###
alias langposix="langset POSIX POSIX"
alias langzh="langset zh_CN.UTF-8 zh_CN:zh"
alias langen="langset en_US.UTF-8 en_US:en"
#### direct2cache head END ####
