# profile.d/81-ls-alias.sh
### direct2cache head START ###
[ -z "$fish_pid" ] && alias ll="ls -lh" && alias la="ls -lAh" && echo aa

alias l=ls
alias l.='ls -d .*'
alias ll.='ll -d .*'
#### direct2cache head END ####
