# profile.d/81-ls-alias.sh
abbr l   ls
abbr ll  'ls -lh'
abbr la  'ls -lAh'
abbr l.  'ls -d .*'
abbr ll. 'll -d .*'
