# profile.d/83-langset-alias.sh
# Used as Fish scripts too.

abbr langposix "langset POSIX POSIX"
abbr langzh    "langset zh_CN.UTF-8 zh_CN:zh"
abbr langen    "langset en_US.UTF-8 en_US:en"
