# profile.d/90-route.sh
#!/bin/sh
### direct2cache head START ###
route()       { command route -n      "$@"; }
traceroute()  { command traceroute -n "$@"; }
#### direct2cache head END ####
