# profile.d/99-aliasfuncs.sh
###### misc ######
path()           { echo "$PATH" | tr : '\n'; }
waitexist()      { waittrue test -e    "$@"; }
