# profile.d/curl-wget.sh
# Used as Fish scripts too.

abbr curltonull 'curl -v -o /dev/null'
abbr wgettonull 'wget -v -O /dev/null'
