# profile.d/df-du.sh
# Used as Fish scripts too.

abbr dfh  'df -h'
abbr dush 'du -sh'
