# profile.d/direnv.sh
if havebin direnv; then
  eval "$(direnv hook bash)"
fi
