# profile.d/editor.sh
# Used as Fish scripts too.

export EDITOR="$editor"
