# profile.d/env.sh
# Used as Fish scripts too.

abbr envclear 'env -i "HOME=$HOME"'
abbr envsort  'env | sort'
abbr envless  'env | sort | less'
