# profile.d/env.sh
abbr envclear 'env -i "HOME=$HOME"'
abbr envsort  'env | sort'
