# profile.d/locale.sh
# Used as Fish scripts too.

[ -n "$LANG" ] && export LC_ALL="$LANG"
