# profile.d/locale.sh
# Used as Fish scripts too.

[ -n "$lang" ] && \
  export LANG="$lang" \
  export LC_ALL="$lang"

[ -n "$language" ] && \
  export LANGUAGE="$language"
