# profile.d/ps-alias.sh
# Used as Fish scripts too.

abbr pse   'ps -e'
abbr psef  'ps -ef'
abbr psme  "ps -u '$USER'"
abbr psmef "ps -u '$USER' -f"
