# profile.d/screen.sh
# Used as Fish scripts too.

abbr scr      'screen -S'
abbr lssc     'screen -ls'
abbr lsscr    'screen -ls'
abbr lsscreen 'screen -ls'
abbr xsc      'screen -x'
abbr xscr     'screen -x'
abbr xscreen  'screen -x'
