# profile.d/screen.sh
# Used as Fish scripts too.

abbr lsscreen 'screen -ls'
abbr xscreen  'screen -x'
