# profile.d/source-pyenv.sh
if havebin pyenv; then
  eval "$(pyenv init -)"
fi
