# profile.d/timestyle.sh
# Used as Fish scripts too.

export TIME_STYLE="$time_style"
