# profile.d/use-lesspipe.sh
# Used as Fish scripts too.

have lesspipe.sh && export LESSOPEN="|lesspipe.sh %s"
