# profile.d/watch.sh
abbr watch      'watch -e -n1 '  # error exit; interval = 1s
abbr watchcolor 'watch -c '
abbr watchdiff  'watchcolor --differences=cumulative '
