# profile.d/watch.sh
alias      watch='watch -e -n1 '  # error exit; interval = 1s
alias watchcolor='watch -c '
alias  watchdiff='watchcolor --differences=cumulative '
